# .config/autostart/annex-tray.desktop
[Desktop Entry]
Type=Application
Name=Annex Tray
Comment=Monitor Git Annex repositories
# Use sh -c to expand $HOME dynamically at runtime
Exec=sh -c "$HOME/.bin/annex-tray"
Icon=folder-remote-symbolic
Terminal=false
X-GNOME-Autostart-enabled=true
